# supervisor/conf.d/web.conf
[program:web]
command=node /app/code/web/server.js
directory=/app/code/web
user=root
autostart=true
autorestart=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
environment=NODE_ENV="production",CONFIG_PATH="/app/data/config.json",DATA_DIR="/app/data",PLAYWRIGHT_BROWSERS_PATH="/app/code/playwright-browsers",SUPERVISOR_SOCK="unix:///run/supervisor.sock"
